# cron.d/renew-ssl

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

#m  h dom mon dow user  command
  0 3 1   */2 *   root  /home/shellicar/docker/shellicar/renew
